# deploy/sub2api.service
[Unit]
Description=GGbond - AI API Gateway Platform
Documentation=https://github.com/pengbin9472/ggbond
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ggbond
Group=ggbond
WorkingDirectory=/opt/ggbond
ExecStart=/opt/ggbond/ggbond
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ggbond

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/ggbond

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
